# bin/init.d/grayson~
#!/bin/sh

#set -x

#
# chkconfig: 35 99 99
# description: Grayson System
#

. /etc/rc.d/init.d/functions

# Requires access to root owned certifiates for https.
GRAYSON_HOME="/opt/grayson/current"
INSTALL_DIR="/opt/grayson/stack"
. $GRAYSON_HOME/../setup.sh

ELEMENTS="rabbitmq httpd celeryd events"

do_start()
{
    for element in $ELEMENTS; do
	/etc/init.d/grayson-${element} start
    done
}
do_stop()
{
    for element in $( echo $ELEMENTS | awk '{ for (i = NF; i > 0; --i) print $i }' ); do
	/etc/init.d/grayson-${element} stop
    done
}

case "$1" in
        start)
                do_start
                ;;
        stop)
                do_stop
                ;;
        restart)
                do_stop
                do_start
                ;;
        *)
                echo "Usage: $0 {start|stop|restart}"
                RETVAL=1
esac

exit $RETVAL
